# com.user.claude-imessage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-imessage</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/USER/.nvm/versions/node/v24.12.0/bin/node</string>
        <string>/Users/USER/Documents/PassiveIncome/SendblueBase/claude-via-imessage/daemon/dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/USER/Documents/PassiveIncome/SendblueBase/claude-via-imessage/daemon</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/USER/.local/log/claude-imessage.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/USER/.local/log/claude-imessage.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/USER/.nvm/versions/node/v24.12.0/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
